# crontab
0 0 * * * cd /app/sticker-repo/s1 && git pull >> /var/log/git-pull.log 2>&1